# house-bot.timer
[Unit]
Description=Run house-bot daily at 9am

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
